# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=crond -fL /dev/stdout

[program:httpd]
command=httpd -u 1000 -c /etc/httpd.conf -f
